# template/sandbox/selfsmith-loop.timer
# selfsmith-loop.timer is the cron replacement: fires selfsmith-loop.service on a schedule.
# The service runs run-iteration.sh, which still takes the single-flight lock and time-box,
# so an overlapping fire is skipped cleanly. Adjust OnCalendar to your domain's cadence
# (how fast does new data actually accumulate between ticks?).

[Unit]
Description=Run the selfsmith self-improvement loop hourly

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
